# profile.d/google-cloud.sh
export CLOUD_SDK_REPO="cloud-sdk-$(lsb_release -cs)"
