# restic/.config/systemd/user/restic-backup.service
[Unit]
Description=restic backup service

[Service]
Type=simple
ExecStart=%h/.local/bin/restic-backup-job

[Install]
WantedBy=default.target
